# service/r53clip.upstart.erb
description "Route53 cluster IP manager"

start on started <%= depend_service %>
stop on stopping <%= depend_service %>

env RUBY_THREAD_MACHINE_STACK_SIZE=1572864

pre-start script
	mkdir -p '<%= run_dir %>'
end script

chdir <%= run_dir %>

exec <%= app_dir %>/bin/foreman run -d <%= app_dir %> -- r53clip start <%= config_path %>

respawn

kill timeout 180

# vim: set ft=upstart ts=2 sw=2 tw=0 :
